# systemd/wispr-flow.service
[Unit]
Description=Wispr Flow (status HUD / dictation) with Hyprland helper
After=graphical-session.target
# No longer depends on wispr-key-remap grab path

[Service]
Type=simple
TimeoutStopSec=5
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
ExecStartPre=/home/wraient/.config/hypr/scripts/wispr-set-caps-ptt.sh
ExecStart=/home/wraient/.config/hypr/scripts/wispr-flow-launch.sh
ExecStartPost=/bin/sleep 2
ExecStartPost=/home/wraient/.config/hypr/scripts/wispr-set-caps-ptt.sh
ExecStartPost=/home/wraient/.config/hypr/scripts/wispr-close-hub-once.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
