# systemd/wispr-key-remap.service
[Unit]
Description=Remap Esc->F20 and Caps->Esc for Wispr Flow (raw EV_KEY)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/hypr/scripts/wispr-key-remap.py
Restart=always
RestartSec=1
# input group / uinput ACL already grant access for this user

[Install]
WantedBy=default.target
