# systemd/wispr-cliphist-cleanup.service.d/override.conf
[Service]
Environment=WISPR_CLIPHIST_DELAY=2
